# supervisord.conf
[supervisord]
logfile = /tmp/supervisord.log
logfile_maxbytes = 50MB
logfile_backups=10
loglevel = info
pidfile = /tmp/supervisord.pid
nodaemon = false
minfds = 1024
minprocs = 200
umask = 022
user = water
identifier = supervisor
nocleanup = true
childlogdir = /tmp
strip_ansi = false

[unix_http_server]
file = /tmp/supervisor.sock
chmod = 0777

[program:kael_web]
command = gunicorn -k gevent manage:app  -b 0.0.0.0:5000 -w 2
autostart=true
autorestart=true
stdout_logfile=/tmp/%(program_name)s.log
stderr_logfile=/tmp/%(program_name)s.log


[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface